# deploy/systemd/anah-autodeploy.service
[Unit]
Description=Anah auto-deploy from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/website
Environment=APP_DIR=/home/ubuntu/website
Environment=BRANCH=main
Environment=PM2_APP_NAME=anah-web
ExecStart=/bin/bash /home/ubuntu/website/scripts/check-and-deploy.sh

